# recipes-support/usb-gadget/files/usb-gadget.service
# /lib/systemd/system/usb-gadget.service — systemd unit for RNDIS gadget

[Unit]
Description=Configure USB RNDIS gadget     
After=local-fs.target                     

[Service]
Type=oneshot                              
ExecStart=/usr/bin/usb-gadget.sh          
RemainAfterExit=yes                        

[Install]
WantedBy=multi-user.target                 

