# recipes-support/usb-gadget/files/usb-gadget.service
# /lib/systemd/system/usb-gadget.service — systemd unit for RNDIS gadget

[Unit]
Description=Configure USB RNDIS gadget      # shown in systemctl status
After=local-fs.target                      # wait until local filesystems are ready

[Service]
Type=oneshot                               # run ExecStart once and exit
ExecStart=/usr/bin/usb-gadget.sh           # the gadget setup script
RemainAfterExit=yes                        # consider the service "active" after running

[Install]
WantedBy=multi-user.target                 # enable at runlevels 3/4/5 (typical)

